# src/rd-frontend/root/etc/services.d/web/run
#!/usr/bin/with-contenv bash

cd /gclient || exit
HOME="/config" exec \
	s6-setuidgid gildedguy node app.js
