# src/alpine/overlay/etc/profile.d/path-fix.sh
#!/bin/env sh

export PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/google-cloud-sdk/bin:/config/.local/bin
